# sni-proxy/DEBIAN/postinst
#!/bin/sh

systemctl daemon-reload
systemctl enable sni-proxy
systemctl start sni-proxy
