# sni-proxy/DEBIAN/prerm
#!/bin/sh
systemctl stop sni-proxy
systemctl disable sni-proxy
systemctl daemon-reload
